# lcr-obmc/meta-lcr/recipes-app/lcr-gpiomonitor-app/files/LCR_GPIO_Mon.service
[Unit]
Description=GPIO Monitor Service
After=network.target

[Service]
ExecStart=/usr/bin/lcrgpiomonitor
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target