# lcr-obmc/meta-lcr/recipes-app/lcr-MANAGER-app/files/LCR_MANAGER.service
[Unit]
Description=BOOT SEQUENCE AND MANAGER
After=systemd-journald.service

[Service]
ExecStart=/usr/bin/MANAGER
Type=simple

[Install]
WantedBy=multi-user.target